# base/bionic/root/etc/services.d/xvfb/run
#!/usr/bin/with-contenv sh

/usr/bin/Xvfb "$DISPLAY" -screen 0 "$SCR_WIDTH"x"$SCR_HEIGHT"x16
