# base/_root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

s6-svwait /var/run/s6/services/novnc


/usr/sbin/nginx
