# base/_root/etc/services.d/nginx/run
#!/bin/bash

set -eu
set -o pipefail

/usr/sbin/nginx
