# base/_root/etc/services.d/vncconfig/run
#!/bin/bash

set -eu
set -o pipefail

s6-setuidgid abc vncconfig -nowin
