# base/_root/etc/services.d/tigervnc/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

exec s6-setuidgid kid /usr/bin/Xvnc "$DISPLAY" -geometry "$SCR_WIDTH"x"$SCR_HEIGHT" -depth 16 -SecurityTypes none -AlwaysShared

