# base/_root/etc/services.d/tigervnc/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

s6-setuidgid abc /usr/bin/Xvnc "$DISPLAY" -geometry "$SCR_WIDTH"x"$SCR_HEIGHT" -depth 16 -PasswordFile=$HOME/.vnc/passwd