# base/bionic/root/etc/services.d/kde/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

export XDG_RUNTIME_DIR="/tmp/runtime-abc"
s6-setuidgid abc /usr/bin/startkde
