# base/bionic/root/etc/services.d/x11vnc/run
#!/usr/bin/with-contenv sh

/usr/bin/x11vnc -display "$DISPLAY" -forever -shared -xrandr -nopw
