# base/focal/root/etc/services.d/kde/run
#!/usr/bin/with-contenv sh

/usr/bin/startplasma-x11
