# base/bionic/root/etc/services.d/novnc/run
#!/usr/bin/with-contenv sh


/novnc/utils/launch.sh --listen 8080 --web /novnc
