# base/_root/etc/services.d/novnc/run
#!/bin/bash

set -eu
set -o pipefail

/novnc/utils/launch.sh
