# base/_root/etc/services.d/novnc/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

s6-svwait /var/run/s6/services/tigervnc


exec s6-setuidgid kid /novnc/utils/launch.sh

