# base/_root/etc/services.d/novnc/run
#!/bin/bash

set -eu
set -o pipefail

s6-setuidgid kid /novnc/utils/launch.sh
